# build/supervisord.d/opensearch-prometheus-exporter.ini
[program:opensearch-prometheus-exporter]
command=python3 main.py
directory=opensearch-prometheus-exporter
user=root
stopsignal=INT
autostart=true
autorestart=true
startsecs=3
startretries=999
stderr_logfile=/var/log/opensearch-prometheus-exporter.err.log
stdout_logfile=/var/log/opensearch-prometheus-exporter.out.log